# scripts/systemd/security-guard-sync.service
[Unit]
Description=Security Guard GCP Sync
After=network.target

[Service]
Type=oneshot
User=gavinfullertx
WorkingDirectory=/opt/security-guard
ExecStart=/usr/bin/podman exec vivint-security-guard python /app/manual_sync.py
StandardOutput=append:/var/log/security-guard/sync.log
StandardError=append:/var/log/security-guard/sync.log
