# scripts/systemd/security-guard-weekly.timer
[Unit]
Description=Weekly VM Maintenance

[Timer]
OnCalendar=Sun *-*-* 20:00:00
Persistent=true

[Install]
WantedBy=timers.target
