# scripts/systemd/security-guard-weekly.service
[Unit]
Description=Weekly Security Guard VM Maintenance
After=network.target

[Service]
Type=oneshot
User=root
ExecStart=/opt/security-guard/scripts/weekly-maintenance.sh
StandardOutput=append:/var/log/security-guard/weekly.log
StandardError=append:/var/log/security-guard/weekly.log
